# airsane-scanner/rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Avahi daemon from S6"

until [ -S /run/dbus/system_bus_socket ]; do
  sleep 1s
done

exec avahi-daemon --no-chroot
